# roles/soe_misc/files/satellite-sync.cron
# THIS FILE IS MANAGED BY ANSIBLE, ANY LOCAL CHANGES WILL BE DESTROYED

@reboot  root  /usr/local/bin/satellite-sync
